# build/shipbuilder.service
[Unit]
Description=Shipbuilder Service

[Service]
Restart=always
EnvironmentFile=/etc/default/shipbuilder
ExecStart=/usr/bin/shipbuilder server

[Install]
WantedBy=multi-user.target
